# _sys/usr/lib/systemd/system/karajo.path
## SPDX-FileCopyrightText: 2023 M. Shulhan <ms@kilabit.info>
## SPDX-License-Identifier: GPL-3.0-or-later

[Unit]
Description=Watch karajo binary

[Path]
PathChanged=/usr/bin/karajo
Unit=systemctl-restart-karajo@%p.service

[Install]
WantedBy=multi-user.target
